# src/scripts/systemd/owfs.service.in
[Unit]
Description=1-wire filesystem FUSE mount
Documentation=man:owfs(1)

[Service]
ExecStart=@bindir@/owfs --server=127.0.0.1 --allow_other /run/owfs
ExecStop=/usr/bin/umount /run/owfs
RuntimeDirectory=owfs

[Install]
WantedBy=multi-user.target
