# src/scripts/systemd/owftpd.service.in
[Unit]
Description=Anonymous FTP server for 1-wire access
Documentation=man:owftpd(1)
Requires=owserver.service
After=owserver.service

[Service]
ExecStart=@bindir@/owftpd --foreground --server=127.0.0.1
User=ow
Group=ow

[Install]
WantedBy=multi-user.target
