# src/scripts/systemd/owhttpd.service.in
[Unit]
Description=Tiny webserver for 1-wire control
Documentation=man:owhttpd(1)
After=avahi-daemon.service

[Service]
ExecStart=@bindir@/owhttpd --foreground --server=127.0.0.1
User=ow
Group=ow

[Install]
WantedBy=multi-user.target
