# src/scripts/systemd/owserver.socket.in
[Unit]
Description=Listening socket for owserver
Documentation=man:owserver(1)

[Socket]
ListenStream=4304

[Install]
WantedBy=sockets.target
