# src/scripts/systemd/owserver.service.in
[Unit]
Description=Backend server for 1-wire control
Documentation=man:owserver(1)

[Service]
Type=notify
NotifyAccess=all
ExecStart=@bindir@/owserver --w1 --foreground
Restart=on-failure
#User=ow
#Group=ow

[Install]
WantedBy=multi-user.target
Also=owserver.socket
